# wucios/releases/noether-forge-v2.4.0/overlay/etc/init.d/wuci-loopback
#!/sbin/openrc-run

description="Bring up the Noether Core loopback interface"

depend() {
  need localmount
  before nftables
}

start() {
  ebegin "Bringing up loopback only"
  ip link set dev lo up
  eend $?
}
